# contrib/init/newyorkcoind.conf
description "NewYorkCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env NEWYORKCOIND_BIN="/usr/bin/newyorkcd"
env NEWYORKCOIND_USER="newyorkcoin"
env NEWYORKCOIND_GROUP="newyorkcoin"
env NEWYORKCOIND_PIDDIR="/var/run/newyorkcd"
# upstart can't handle variables constructed with other variables
env NEWYORKCOIND_PIDFILE="/var/run/newyorkcd/newyorkcd.pid"
env NEWYORKCOIND_CONFIGFILE="/etc/newyorkcoin/newyorkcoin.conf"
env NEWYORKCOIND_DATADIR="/var/lib/newyorkcd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # newyorkcd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$NEWYORKCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run newyorkcd."
        echo "The setting must appear in $NEWYORKCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"NewYorkCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$NEWYORKCOIND_PIDDIR"
    chmod 0755 "$NEWYORKCOIND_PIDDIR"
    chown $NEWYORKCOIND_USER:$NEWYORKCOIND_GROUP "$NEWYORKCOIND_PIDDIR"
    chown $NEWYORKCOIND_USER:$NEWYORKCOIND_GROUP "$NEWYORKCOIND_CONFIGFILE"
    chmod 0660 "$NEWYORKCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$NEWYORKCOIND_PIDFILE" \
    --chuid $NEWYORKCOIND_USER:$NEWYORKCOIND_GROUP \
    --exec "$NEWYORKCOIND_BIN" \
    -- \
    -pid="$NEWYORKCOIND_PIDFILE" \
    -conf="$NEWYORKCOIND_CONFIGFILE" \
    -datadir="$NEWYORKCOIND_DATADIR" \
    -disablewallet \
    -daemon

